# assemble/mysql/5.7.16-i686/service/mysqld
#!/bin/bash
#
# chkconfig: 3 20 80
# description: Starts and stops the mysqld server

# Source function library.
if [ -f /etc/init.d/functions ] ; then
    . /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ] ; then
    . /etc/rc.d/init.d/functions
else
    exit 1
fi

# Avoid using root's TMPDIR
unset TMPDIR

RETVAL=0
BINFILE=/usr/local/mysql/bin/mysqld_safe
PIDFILE={{sBinPidInstallFile}}
ADMINFILE=/usr/local/mysql/bin/mysqladmin
CFGFILE=/usr/local/mysql/my.cnf
KIND=mysqld

start() {
    /bin/echo -n $"Starting $KIND services: "
    $BINFILE {{sBinStartParam}} &
    RETVAL=$?
    success
    /bin/echo
    [ $RETVAL -eq 0 ]
        RETVAL=1
    return $RETVAL
}

stop() {
    /bin/echo -n $"Shutting down $KIND services: "
    $ADMINFILE -u root -p shutdown
    RETVAL=$?
    success
    /bin/echo
    [ $RETVAL -eq 0 ] && rm -f $PIDFILE 
    return $RETVAL
}

restart() {
    stop
    start
}

status() {
    if [ ! -r $PIDFILE ]; then
        /bin/echo "mysqld is stopped"
        exit 0
    fi

    PID=`/bin/cat $PIDFILE`
    if ps -p $PID | grep -q $PID; then
        /bin/echo "mysqld (pid $PID) is running..."
    else
        /bin/echo "mysqld dead but pid file exists"
    fi
}

# Check that we can write to it... so non-root users stop here
[ -w $BINFILE ] || exit 4

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    status)
        status
        ;;
    *)
        /bin/echo $"Usage: $0 {start|stop|restart|status}"
        exit 2
esac

exit $?